# proxy-config/systemd/tor.service.d/35_securedrop-proxy-config.conf
[Unit]
# Won't compete with 30_qubes.conf. Both ConditionPathExists should be applied
ConditionPathExists=/var/run/qubes-service/securedrop-proxy-config
