# log/systemd/rsyslog.service.d/35_securedrop-xen.conf
# Our plugin invokes qrexec-client-vm, which needs /dev/xen to
# survive the PrivateDevices=yes hardening
[Service]
BindPaths=-/dev/xen
